# config/runit/log/run
#!/usr/bin/env bash

svlogd -tt {env.service_dir}/log 2>&1
